# systemd/brian-cleanup.timer
[Unit]
Description=Package repository manager, periodic cleanup

[Timer]
OnCalendar=Sat *-*-* 6:05
Persistent=true

[Install]
WantedBy=timers.target
